# timer/check-aurto-git-trigger.service
[Unit]
Description=Lazily sets to check for aurto git updates

[Service]
Type=oneshot
ExecStart=/usr/bin/touch /usr/lib/aurto/check-vcs
