# timer/update-aurto.timer
[Unit]
Description=Updates 'aurto' repo every hour

[Timer]
OnCalendar=*-*-* *:30:00
OnActiveSec=5s

[Install]
WantedBy=timers.target
