# timer/update-aurto.service
[Unit]
Description=Updates 'aurto' repo

[Service]
Type=oneshot
ExecStart=/usr/lib/aurto/update-aurto
