# timer/update-aurto.service
[Unit]
Description=Updates aurto repos
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/lib/aurto/update-aurto
