# timer/update-aurto-startup.timer
[Unit]
Description=Updates aurto repos on startup

[Timer]
OnActiveSec=5s
Unit=update-aurto.service

[Install]
WantedBy=timers.target
