# timer/update-aurto.service
[Unit]
Description=Updates 'aurto' repo

[Service]
ExecStart=/usr/lib/aurto/update-aurto
