# timer/update-aurto.timer
[Unit]
Description=Updates aurto repos every hour

[Timer]
OnCalendar=*-*-* *:25:00
RandomizedDelaySec=10min

[Install]
WantedBy=timers.target
